# deploy/skaelia-tunnel.service
[Unit]
Description=Skaelia Prospection - tunnel Cloudflare (croisia.me)
After=network-online.target skaelia-app.service
Wants=network-online.target

[Service]
Type=simple
User=skaelia
# Le tunnel envoie le trafic de croisia.me vers l'app locale (127.0.0.1:5173).
ExecStart=/usr/local/bin/cloudflared tunnel --config /opt/skaelia/.cloudflared/config.yml run skaelia-prospection
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
